# btcpayserver.service
[Unit]
Description=btcpayserver
Requires=nbxplorer.service
Requires=lnd.service
After=nbxplorer.service
After=lnd.service

[Service]
User=elmeri
WorkingDirectory=/home/elmeri/ssd-fullnode
ExecStart=btcpayserver -d /home/elmeri/ssd-fullnode/bitcoin-blockchain-datadir/btcpayserver
Restart=always
RestartSec=5
StartLimitInterval=400
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
